# lightsail/gt-session-refresh.timer
[Unit]
Description=Run Gran Turismo session refresh every 15 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min
Unit=gt-session-refresh.service
Persistent=true

[Install]
WantedBy=timers.target
